# deploy/systemd/gira-station.service
[Unit]
Description=GIRA station snapshot scraper
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=gira
Group=gira
WorkingDirectory=/opt/gira-watcher/packages/scraper
EnvironmentFile=/opt/gira-watcher/.env
ExecStart=/opt/gira-watcher/packages/scraper/.venv/bin/python -m scraper.main --run-type station
TimeoutStartSec=240

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/gira-watcher
